# website/backend/config/send-sonde-email.timer
[Unit]
Description=Hourly send_sonde_email run

[Timer]
OnCalendar=*-*-* *:30:00
Persistent=true

[Install]
WantedBy=timers.target
